# keyoxide-web.service
[Unit]
Description=Keyoxide (Your online passport)
After=syslog.target
After=network.target

[Service]
User=keyoxide
Group=www-data
WorkingDirectory=/opt/keyoxide-web/
ExecStart=/usr/bin/node /opt/keyoxide-web/src/index.js
Restart=always
RestartSec=2s
Environment=PORT=5000

[Install]
WantedBy=multi-user.target
